# .local/services/nginx
#!/sbin/openrc-run

name="${RC_SVCNAME}"
supervisor=supervise-daemon
USER=jan
GROUP=jan
COMPOSE_FILE="/var/lib/${RC_SVCNAME}/compose.yaml"

start() {

	ebegin "Starting ${RC_SVCNAME}"
	# Stop podman-compose down from possible previous stop() invocation
	${supervisor} "${RC_SVCNAME}" \
		--stop podman-compose
	# Start podman-compose up
	${supervisor} "${RC_SVCNAME}" \
		--stdout "/var/log/${RC_SVCNAME}/output-up.out" \
		--stderr "/var/log/${RC_SVCNAME}/output-up.err" \
		--user "${USER}" --group "${GROUP}" \
		--start podman-compose -- -f "${COMPOSE_FILE}" up --force-recreate
	eend $?

}

stop() {

	ebegin "Stopping ${RC_SVCNAME}"
	# Stop podman-compose up from previous start() invocation
	${supervisor} "${RC_SVCNAME}" \
		--stop podman-compose
	# Start podman-compose down
	${supervisor} "${RC_SVCNAME}" \
		--stdout "/var/log/${RC_SVCNAME}/output-down.out" \
		--stderr "/var/log/${RC_SVCNAME}/output-down.err" \
		--user "${USER}" --group "${GROUP}" \
		--start podman-compose -- -f "${COMPOSE_FILE}" down -t 0
	eend $?

}

