# podman/jellyfin/service
#!/sbin/openrc-run

name="${RC_SVCNAME}"
supervisor=supervise-daemon
USER=jellyfin
GROUP=jellyfin
COMPOSE_FILE="/var/lib/${RC_SVCNAME}/compose.yaml"

CHECK_NVIDIA_PATH="/dev/nvidia-modeset"
DEVICES_NVIDIA="nvidia.com/gpu=all"
DEVICE_PATH_CARD_NONVIDIA="/dev/dri/by-path/pci-0000:00:02.0-card"
DEVICE_PATH_RENDER_NONVIDIA="/dev/dri/by-path/pci-0000:00:02.0-render"
# DEVICE_PATH_CARD_NVIDIA="/dev/dri/by-path/pci-0000:01:00.0-card"
# DEVICE_PATH_RENDER_NVIDIA="/dev/dri/by-path/pci-0000:01:00.0-render"
DEVICE_PATH_CARD_TARGET="/dev/dri/card0"
DEVICE_PATH_RENDER_TARGET="/dev/dri/renderD128"
CONTAINER_NAME_NONVIDIA=jellyfin
CONTAINER_NAME_NVIDIA=jellyfin_nvidia

export CONTAINER_NAME
export DEVICES
export USER_ID
export GROUP_ID

start_pre() {

	# get UID and GID
	USER_ID=$(id -u ${USER})
	GROUP_ID=$(id -g ${GROUP})

	# set container name and devices
	if ls "${CHECK_NVIDIA_PATH}" > /dev/null 2>&1
	then
		CONTAINER_NAME=${CONTAINER_NAME_NVIDIA}
		DEVICES="$DEVICES_NVIDIA"
	else
		CONTAINER_NAME=${CONTAINER_NAME_NONVIDIA}
		DEVICE_PATH_CARD_NONVIDIA=$(readlink -f ${DEVICE_PATH_CARD_NONVIDIA})
		DEVICE_PATH_RENDER_NONVIDIA=$(readlink -f ${DEVICE_PATH_RENDER_NONVIDIA})
		DEVICES="${DEVICE_PATH_CARD_NONVIDIA},${DEVICE_PATH_CARD_TARGET}:${DEVICE_PATH_RENDER_NONVIDIA},${DEVICE_PATH_RENDER_TARGET}"
	fi
	return 0

}

start() {

	ebegin "Starting ${RC_SVCNAME}"
	# Stop podman-compose down from possible previous stop() invocation
	${supervisor} "${RC_SVCNAME}" \
		--stop podman-compose
	# Start podman-compose up
	${supervisor} "${RC_SVCNAME}" \
		--stdout "/var/log/${RC_SVCNAME}/output-up.out" \
		--stderr "/var/log/${RC_SVCNAME}/output-up.err" \
		--user "${USER}" --group "${GROUP}" \
		--start podman-compose -- -f "${COMPOSE_FILE}" up
	eend $?

}

stop_pre() {
	start_pre
}

stop() {

	ebegin "Stopping ${RC_SVCNAME}"
	# Stop podman-compose up from previous start() invocation
	${supervisor} "${RC_SVCNAME}" \
		--stop podman-compose
	# Start podman-compose down
	${supervisor} "${RC_SVCNAME}" \
		--stdout "/var/log/${RC_SVCNAME}/output-down.out" \
		--stderr "/var/log/${RC_SVCNAME}/output-down.err" \
		--user "${USER}" --group "${GROUP}" \
		--start podman-compose -- -f "${COMPOSE_FILE}" down -t 0
	eend $?
	
}
